# config/docker/runit/echo/run
#!/bin/bash
cd /app
exec bundle exec thin start -S /tmp/sockets/thin.sock -e production
